# duckdns/duckdns.service
[Unit]
Description=Duck DNS Updater

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/duck.sh >/dev/null 2>&1


