# borg_system/borg-system.service
[Unit]
Description=Backup system to rsync.net

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/borg_system
